# templates/default/etc/systemd/system/karaf.service.erb
[Unit]
SourcePath=/etc/init.d/karaf
Description=karaf
After=remote-fs.target systemd-journald-dev-log.socket

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/etc/init.d/karaf start
ExecStop=/etc/init.d/karaf stop

[Install]
WantedBy=multi-user.target
