# templates/default/etc/systemd/system/node-ff.service.erb
[Unit]
SourcePath=/etc/init.d/node-ff
Description=node-ff
After=remote-fs.target systemd-journald-dev-log.socket

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/etc/init.d/node-ff start
ExecStop=/etc/init.d/node-ff stop

[Install]
WantedBy=multi-user.target
