# templates/default/etc/systemd/system/node-chrome.service.erb
[Unit]
SourcePath=/etc/init.d/node-chrome
Description=node-chrome
After=remote-fs.target systemd-journald-dev-log.socket

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/etc/init.d/node-chrome start
ExecStop=/etc/init.d/node-chrome stop

[Install]
WantedBy=multi-user.target
